# src/examples/etc/service/thingpin/run
#!/bin/bash

exec thingpin -c /etc/thingpin/thingpin-config.yml run 2>&1 | logger -t thingpin
