# src/examples/etc/service/thingpin/run
#!/bin/bash

exec thingpin \
     -p /var/run/thingpin.pid \
     -c /etc/thingpin/thingpin-config.yml \
     run 2>&1 | logger -t thingpin
